# hassnowplaying.service
[Unit]
Description=Home Assistant now playing

[Service]
Type=notify
ExecStart=/usr/sbin/testapp

[Install]
WantedBy=multi-user.target